# azure/pg-grafana/stack/supervisord.conf
; Programs only — the base /etc/supervisor/supervisord.conf provides the
; [supervisord] section and includes this file from conf.d/.

[program:mimir]
command=/usr/local/bin/mimir -config.file=/etc/mimir/mimir.yaml -target=all
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:alloy]
command=/usr/bin/alloy run /etc/alloy/config.alloy --storage.path=/var/lib/alloy --server.http.listen-addr=0.0.0.0:12345
autorestart=true
priority=20
environment=PG_EXPORTER_DSN="%(ENV_PG_EXPORTER_DSN)s"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:grafana]
command=/usr/share/grafana/bin/grafana server --homepath=/usr/share/grafana --config=/etc/grafana/grafana.ini
directory=/usr/share/grafana
autorestart=true
priority=30
environment=GF_PATHS_PROVISIONING="/etc/grafana/provisioning",GF_SECURITY_ADMIN_USER="admin",GF_SECURITY_ADMIN_PASSWORD="admin",GF_USERS_DEFAULT_THEME="light"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
